# templates/pushbot-gtalk.conf.j2
respawn

exec su - {{ pushbot_username }} -c 'cd {{ pushbot_root }} ; bin/pushbot-gtalk >> {{ pushbot_log }}/pushbot.gtalk.log 2>&1'

start on runlevel [2345]

stop on runlevel [016]
